# linuxconf/dotfiles.timer
[Unit]
Description=Run dotfile.service every 60 minutes

[Timer]
OnCalendar=*:0/59

