# linuxconf/dotfiles.timer
[Unit]
Description=Run dotfile.service every 60 minutes

[Timer]
OnBootSec=1min
OnCalendar=daily
Persistent=true

